# deploy/traffic-go.service
[Unit]
Description=traffic-go network flow monitor
Documentation=man:systemd.service(5)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_GROUP__
WorkingDirectory=__INSTALL_ROOT__
ExecStartPre=-__MODPROBE_BIN__ nf_conntrack
ExecStartPre=-__SYSCTL_BIN__ -w net.netfilter.nf_conntrack_acct=1
ExecStart=__INSTALL_BIN__ -config __INSTALL_CONFIG__
Restart=on-failure
RestartSec=3
TimeoutStartSec=30
LimitNOFILE=65536
UMask=0007
NoNewPrivileges=true
AmbientCapabilities=CAP_NET_ADMIN CAP_SYS_PTRACE CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_NET_ADMIN CAP_SYS_PTRACE CAP_DAC_READ_SEARCH

[Install]
WantedBy=multi-user.target
